# branches/mgar/minimalsmf/trunk/files/cswminimalsmf.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='system/cswminimalsmf' type='service' version='0'>
    <create_default_instance enabled='false'/>
    <single_instance/>
    <exec_method name='start' type='method' exec='/opt/csw/lib/svc/method/cswminimalsmf' timeout_seconds='600'>
      <method_context/>
    </exec_method>
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'>
      <method_context/>
    </exec_method>
    <exec_method name='refresh' type='method' exec=':kill -HUP' timeout_seconds='60'>
      <method_context/>
    </exec_method>
    <stability value='Unstable'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>Just sleeping</loctext>
      </common_name>
      <documentation>
        <manpage title='bash' section='1' manpath='/usr/share/man'/>
      </documentation>
    </template>
  </service>
</service_bundle>
